# systemd/headphones-plug-controller.service.in
[Unit]
Description=Headphones plug detector and playback controller
After=pulseaudio.service dbus.socket

[Service]
ExecStart=@bindir@/headphones-plug-controller
Restart=on-failure

[Install]
WantedBy=default.target
